# debian/chd2iso-fuse/lib/systemd/system/mnt-retronas-roms-sony-playstation2-iso.automount
[Unit]
Description=Automount CHD→ISO PS2

[Automount]
Where=/mnt/retronas/roms/sony/playstation2/iso

[Install]
WantedBy=multi-user.target
